# init/hybrid-gfx.conf
start on (starting lightdm
          or starting kdm
          or starting mdm
          or starting xdm
          or starting lxdm)
task
exec hybrid-detect
